# deploy/bitnet.service

[Unit]
Description=BitNet inference service (llama-cli)
After=network.target


[Service]
Type=simple
User=bitnet
Group=bitnet
WorkingDirectory=/workspaces/BitNet
# Load optional environment file (create /etc/bitnet/bitnet.env with VAR=VAL lines)
EnvironmentFile=-/etc/bitnet/bitnet.env
Environment=MODEL_PATH=/workspaces/BitNet/models/BitNet-b1.58-2B-4T/ggml-model-i2_s.gguf
Environment=THREADS=4
Environment=PROMPT="You are a helpful assistant."

# Resource limits (adjust per profile: small/medium/large)
MemoryMax=%MEMORY_MAX%
CPUQuota=%CPU_QUOTA%
LimitNOFILE=%LIMIT_NOFILE%

# Prefer journald output; if LOG_TO_JOURNAL is set to 0 in env file, logs will be redirected to /var/log/bitnet/bitnet.log
Environment=LOG_TO_JOURNAL=1

ExecStart=/usr/bin/env bash -lc 'if [ "${LOG_TO_JOURNAL}" = "1" ]; then exec /workspaces/BitNet/scripts/run_bitnet.sh "$MODEL_PATH" "$PROMPT" "$THREADS" 256; else exec /workspaces/BitNet/scripts/run_bitnet.sh "$MODEL_PATH" "$PROMPT" "$THREADS" 256 >> /var/log/bitnet/bitnet.log 2>&1; fi'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
